# 93-module-git/ansible_dir/myapp.service.j2
[Unit]
Description= Flask myapp
After=network.target

[Service]
User=root
ExecStart=/opt/myapp/sample_api.py

[Install]
WantedBy=multi-user.target